# package/classifier/DEBIAN/postrm
#!/bin/sh
update-rc.d -f classifier remove
